# supervisor-app.conf
[program:app-uwsgi]
command=/usr/local/bin/uwsgi --ini /code/uwsgi.ini
autostart = true
autorestart = true
startretries = 3
stdout_logfile = /code/logs/supervisord.log